# files/openhab.conf
[program:openhab]
command=/usr/local/bin/openhab.sh
directory=/opt/openhab
stopwaitsecs=30
